# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-core/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start MariaDB service
# ==============================================================================

# Start mariadb
bashio::log.info "Starting MariaDB"
mkdir -p /run/mysqld
exec mysqld --datadir="${MARIADB_DATA}" --user=root < /dev/null

echo "[galera]" >> /etc/mysql/my.cnf
echo "wsrep_cluster_name=${galera_cluster_name}" >> /etc/mysql/my.cnf
echo "wsrep_cluster_address=${galera_cluster_address}" >> /etc/mysql/my.cnf
echo "wsrep_node_address=${galera_node_address}" >> /etc/mysql/my.cnf
echo "wsrep_sst_method=${galera_sst_method}" >> /etc/mysql/my.cnf
